# chef.sv.conf
[program:Chef]
command=/bin/bash -c "/home/hanlon/Chef/entrypoint.sh"
redirect_stderr=true
